# deploy/systemd/qwen-vllm.service
[Unit]
Description=Qwen vLLM OpenAI API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/qwen-web
EnvironmentFile=/etc/qwen-web/vllm.env
ExecStartPre=/bin/bash -lc 'test -f /etc/qwen-web/vllm.env'
ExecStart=/bin/bash -lc '\
  PYTHON_BIN_RESOLVED="$${PYTHON_BIN:-python3}"; \
  PYTHON_BIN_DIR="$$(dirname "$$PYTHON_BIN_RESOLVED")"; \
  export PATH="$$PYTHON_BIN_DIR:$$PATH"; \
  ARGS=( \
    --model "$$MODEL_PATH" \
    --served-model-name "$${SERVED_MODEL_NAME:-qwen}" \
    --trust-remote-code \
    --host "$${VLLM_HOST:-0.0.0.0}" \
    --port "$${VLLM_PORT:-8000}" \
    --gpu-memory-utilization "$${GPU_MEMORY_UTILIZATION:-0.85}" \
    --max-model-len "$${MAX_MODEL_LEN:-65536}" \
    --kv-cache-dtype "$${KV_CACHE_DTYPE:-fp8}" \
    --max-num-batched-tokens "$${MAX_NUM_BATCHED_TOKENS:-4096}" \
    --enable-prefix-caching \
  ); \
  if [[ -n "$${LIMIT_MM_PER_PROMPT:-}" ]]; then \
    ARGS+=(--limit-mm-per-prompt "$$LIMIT_MM_PER_PROMPT"); \
  fi; \
  if [[ -n "$${EXTRA_VLLM_ARGS:-}" ]]; then \
    read -r -a EXTRA_ARGS <<< "$$EXTRA_VLLM_ARGS"; \
    ARGS+=("$${EXTRA_ARGS[@]}"); \
  fi; \
  exec "$$PYTHON_BIN_RESOLVED" -m vllm.entrypoints.openai.api_server "$${ARGS[@]}" \
'
Restart=always
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=60
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
